# assets/PortableGit/etc/profile.d/env.sh
# Add bin path in the home directory ontop of the PATH variable
export PATH="$HOME/bin:$PATH"
